# clippy.service
[Unit]
Description=Clippy — local LLM chat assistant (Flask/Ollama)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=coreconduit
WorkingDirectory=/home/coreconduit/rpi-clippy
ExecStart=/usr/bin/python3 /home/coreconduit/rpi-clippy/app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clippy

[Install]
WantedBy=multi-user.target
